# monitor.conf
[program:monitor]
user=root
directory=/home/pi/Src/BagMonitor
command=/home/pi/Src/BagMonitor/monitorBoard.pl
autostart=true
autorestart=true
